# services/fix-displays.service
[Unit]
Description=Fix dual display setup
After=graphical.target
Before=touch-ui.service touch-response.service

[Service]
Type=oneshot
User=brother
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/brother/.Xauthority
ExecStart=/home/brother/fix-displays.sh
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=graphical.target