# services/touch-companion.service
[Unit]
Description=Touch Companion Service
After=network.target

[Service]
User=pi
WorkingDirectory=/home/pi/touch-companion
ExecStart=/home/pi/touch-companion/.venv/bin/python /home/pi/touch-companion/main.py
Restart=on-failure
RestartSec=5
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=touch-companion
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target